# contrib/systemd/acts.service
[Unit]
Description=acts Tarsnap backup

[Service]
Type=simple
ExecStart=/usr/local/bin/acts
